# deploy/backup/today-list-mysql-backup.service
[Unit]
Description=Back up the Today List MySQL database
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/sbin/today-list-mysql-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/backups/today-list/mysql /run/lock
